# cron.d/FireMotD
00 00 * * * root perl -e 'sleep int(rand(3600))' && /usr/local/bin/FireMotD -S &>/dev/null
