# cron.d/FireMotD
15 18 * * * root /usr/local/bin/FireMotD -S
